# supervisord.ini
[supervisord]
nodaemon = true
user = root
logfile = logs/supervisor.log
pidfile = pids/supervisor.pid

[program:daphne]
command = daphne backend.asgi:application -b 0.0.0.0 --access-log logs/daphne-access.log

[program:background-tasks]
command = python manage.py process_tasks
