# supervisord.ini
[supervisord]
nodaemon = true
user = root
logfile = logs/supervisor.log
pidfile = pids/supervisor.pid

[program:gunicorn]
command = gunicorn backend.wsgi --access-logfile logs/gunicorn-access.log --log-file logs/gunicorn.log -b 0.0.0.0

[program:background-tasks]
command = python manage.py process_tasks
